# systemd/nadirclaw.service
[Unit]
Description=NadirClaw — Intelligent LLM Router Proxy
Documentation=https://github.com/doramirdor/NadirClaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clawuser
Group=clawuser
WorkingDirectory=/home/clawuser

# NadirClaw reads model config and API keys from environment variables.
# OPENROUTER_API_KEY is forwarded to LiteLLM for all OpenRouter model calls.
EnvironmentFile=/home/clawuser/personal-os/.env

ExecStart=/usr/local/bin/nadirclaw serve
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nadirclaw

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/clawuser/.nadirclaw

[Install]
WantedBy=multi-user.target
