# systemd/personal-os-heartbeat.timer
[Unit]
Description=Personal AI OS — Heartbeat Timer (every 55 minutes)
Requires=personal-os-heartbeat.service

[Timer]
# Run 55 minutes after boot, then every 55 minutes
OnBootSec=55min
OnUnitActiveSec=55min
Unit=personal-os-heartbeat.service
AccuracySec=1min

[Install]
WantedBy=timers.target
