# systemd/personal-os-heartbeat.service
[Unit]
Description=Personal AI OS — Agent Heartbeat (one-shot)
After=network-online.target personal-os.service

[Service]
Type=oneshot
User=clawuser
Group=clawuser
WorkingDirectory=/home/clawuser/personal-os
EnvironmentFile=/home/clawuser/personal-os/.env

# Trigger heartbeat for all agents via the OpenClaw CLI
ExecStart=/usr/bin/env openclaw heartbeat --all --config /home/clawuser/.openclaw/openclaw.json

StandardOutput=journal
StandardError=journal
SyslogIdentifier=personal-os-heartbeat
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/clawuser/.openclaw
